# examples/systemd/meural-mcp-api.service
[Unit]
Description=MeuralMCP localhost REST API and remote MCP endpoint
Wants=network-online.target
After=network-online.target meural-mcp-daemon.service

[Service]
Type=simple
User=meural-mcp
Group=meural-mcp
Environment=MEURAL_MCP_STORAGE_DIR=/var/lib/meural-mcp
ExecStart=/opt/meural-mcp/.venv/bin/meural-mcp --storage-dir /var/lib/meural-mcp serve --host 127.0.0.1 --port 8733
Restart=always
RestartSec=15
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ReadWritePaths=/var/lib/meural-mcp

[Install]
WantedBy=multi-user.target
